# LaunchAgents/com.household.checkin.evening.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.household.checkin.evening</string>
    <key>ProgramArguments</key>
    <array>
        <string>{{VAULT_ROOT}}/scripts/household-runner</string>
        <string>run-automation.sh</string>
        <string>checkin</string>
        <string>Evening</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>20</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>{{VAULT_ROOT}}/logs/com.household.checkin.evening.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>{{VAULT_ROOT}}/logs/com.household.checkin.evening.stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin:{{GEMINI_BIN_DIR}}</string>
        <key>VAULT_ROOT</key>
        <string>{{VAULT_ROOT}}</string>
    </dict>
</dict>
</plist>